# bingwall.timer
[Unit]
Description=Run bingwall daily

[Timer]
OnCalendar=daily
Persistent=true

